# base-system/usrroot/etc/systemd/system/hsync.service
[Unit]
Description=huronOS sync service

[Service]
User=root
Group=root
ExecStart=/usr/sbin/hsync.downloader
Restart=always

[Install]
WantedBy=multi-user.target
