# base-system/usrroot/usr/lib/systemd/system/hsync.service
[Unit]
Description=huronOS sync service
Requires=network.target
After=network.target
After=multi-user.target

[Service]
User=root
Group=root
ExecStart=/usr/lib/hsync/hsync.sh --routine-sync
Restart=on-abnormal
SuccessExitStatus=0